# deploy/weather-backend.service
[Unit]
Description=Weather App Backend Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/weather-backend/current
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
EnvironmentFile=/var/www/weather-backend/shared/backend.env
Environment=PORT=3002
# Resource Limits (1GB Server shared with Library App)
MemoryLimit=250M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
